# image/stageB/03-app/files/home/pi/beerpro/scripts/beerpro.service
[Unit]
Description=Beer Pro Scoreboard & Replay System
After=network.target sound.target uap0.service
Wants=network.target

[Service]
Type=notify
User=pi
WorkingDirectory=/home/pi/beerpro
ExecStart=/home/pi/beerpro/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beerpro

# Watchdog — if the app stops notifying for 30s, systemd kills and restarts it
WatchdogSec=30
NotifyAccess=main

# Give the network stack a moment to settle before starting
ExecStartPre=/bin/sleep 3

# Allow the pi user to restart this service without a password (for OTA updates)
# (sudoers entry is configured in the image build)

[Install]
WantedBy=multi-user.target
